# supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=/usr/bin/dotnet /app/aspnetapp.dll

[program:ssh]
command=/usr/sbin/sshd -D